# updater.crontab
*/15 * * * * /srv/Data-Science/run.sh >> /srv/Data-Science/logs/cron.log 2>&1
*/15 * * * * /srv/Data-Science/addusers.sh >> /srv/Data-Science/logs/cron.log 2>&1
*/15 * * * * echo "called at:" >> /srv/Data-Science/logs/cron.log 2>&1
*/15 * * * * date >> /srv/Data-Science/logs/cron.log 2>&1
* * 1 * * echo "run and addusers start logs:" > /srv/Data-Science/logs/cron.log
18,48 * * * * /srv/Data-Science/updateaverages.sh >> /srv/Data-Science/logs/cron.log 2>&1
18,48 * * * * echo "called at:" >> /srv/Data-Science/logs/cron.log 2>&1
18,48 * * * * date >> /srv/Data-Science/logs/cron.log 2>&1
